# root/etc/s6-overlay/s6-rc.d/svc-orca-watchdog/run
#!/usr/bin/with-contenv bash

# OrcaSlicer Watchdog Service
# Monitors OrcaSlicer for hangs and memory issues, performs automatic recovery

WATCHDOG_INTERVAL=30  # Check every 30 seconds
MEMORY_THRESHOLD=80   # Memory usage percentage threshold
CPU_HANG_THRESHOLD=5  # Minutes of 100% CPU = potential hang
RESPONSE_TIMEOUT=10   # Seconds to wait for process response

log() {
    echo "[$(date '+%Y-%m-%d %H:%M:%S')] WATCHDOG: $1"
}

get_orca_pid() {
    pgrep -f "orca-slicer|OrcaSlicer|orca_slicer" | head -1
}

get_memory_usage() {
    local pid=$1
    if [ -n "$pid" ] && [ -f "/proc/$pid/status" ]; then
        local vmrss=$(grep VmRSS /proc/$pid/status 2>/dev/null | awk '{print $2}')
        echo "${vmrss:-0}"
    else
        echo "0"
    fi
}

check_process_responsive() {
    local pid=$1
    if [ -n "$pid" ]; then
        # Check if process responds to signals
        kill -0 "$pid" 2>/dev/null
        return $?
    fi
    return 1
}

clear_gpu_cache() {
    log "Clearing GPU caches..."
    rm -rf /config/.cache/mesa_shader_cache/* 2>/dev/null
    rm -rf /tmp/mesa* 2>/dev/null
    sync
}

restart_orca() {
    log "Initiating OrcaSlicer restart..."
    local pid=$(get_orca_pid)
    
    if [ -n "$pid" ]; then
        # Try graceful shutdown first
        log "Sending SIGTERM to PID $pid"
        kill -TERM "$pid" 2>/dev/null
        sleep 5
        
        # Force kill if still running
        if kill -0 "$pid" 2>/dev/null; then
            log "Force killing PID $pid"
            kill -9 "$pid" 2>/dev/null
        fi
    fi
    
    # Clear caches
    clear_gpu_cache
    
    # Trigger memory cleanup
    sync
    echo 3 > /proc/sys/vm/drop_caches 2>/dev/null || true
    
    log "OrcaSlicer restart completed, waiting for autostart..."
}

log "Starting OrcaSlicer watchdog service"
log "Memory threshold: ${MEMORY_THRESHOLD}%"
log "Check interval: ${WATCHDOG_INTERVAL}s"

HIGH_CPU_COUNT=0

while true; do
    sleep "$WATCHDOG_INTERVAL"
    
    ORCA_PID=$(get_orca_pid)
    
    if [ -z "$ORCA_PID" ]; then
        # OrcaSlicer not running, reset counters
        HIGH_CPU_COUNT=0
        continue
    fi
    
    # Check if process is responsive
    if ! check_process_responsive "$ORCA_PID"; then
        log "OrcaSlicer (PID $ORCA_PID) is not responsive"
        restart_orca
        HIGH_CPU_COUNT=0
        continue
    fi
    
    # Check memory usage (in KB)
    MEM_USAGE=$(get_memory_usage "$ORCA_PID")
    MEM_USAGE_MB=$((MEM_USAGE / 1024))
    
    # Get total system memory
    TOTAL_MEM=$(grep MemTotal /proc/meminfo | awk '{print $2}')
    MEM_PERCENT=$((MEM_USAGE * 100 / TOTAL_MEM))
    
    if [ "$MEM_PERCENT" -gt "$MEMORY_THRESHOLD" ]; then
        log "High memory usage detected: ${MEM_USAGE_MB}MB (${MEM_PERCENT}%)"
        log "Clearing shader cache to free memory..."
        clear_gpu_cache
        
        # If still high after cache clear, restart
        sleep 5
        NEW_MEM=$(get_memory_usage "$ORCA_PID")
        NEW_MEM_PERCENT=$((NEW_MEM * 100 / TOTAL_MEM))
        
        if [ "$NEW_MEM_PERCENT" -gt "$MEMORY_THRESHOLD" ]; then
            log "Memory still high after cache clear, restarting OrcaSlicer"
            restart_orca
        fi
        HIGH_CPU_COUNT=0
        continue
    fi
    
    # Check for CPU hang (stuck at 100%)
    CPU_USAGE=$(ps -p "$ORCA_PID" -o %cpu= 2>/dev/null | cut -d. -f1)
    if [ "${CPU_USAGE:-0}" -gt 95 ]; then
        HIGH_CPU_COUNT=$((HIGH_CPU_COUNT + 1))
        log "High CPU detected: ${CPU_USAGE}% (count: $HIGH_CPU_COUNT)"
        
        if [ "$HIGH_CPU_COUNT" -ge $((CPU_HANG_THRESHOLD * 60 / WATCHDOG_INTERVAL)) ]; then
            log "CPU hang detected, restarting OrcaSlicer"
            restart_orca
            HIGH_CPU_COUNT=0
        fi
    else
        HIGH_CPU_COUNT=0
    fi
    
    # Periodic cache maintenance (every 10 minutes)
    if [ $((SECONDS % 600)) -lt "$WATCHDOG_INTERVAL" ]; then
        CACHE_SIZE=$(du -sm /config/.cache/mesa_shader_cache 2>/dev/null | cut -f1)
        if [ "${CACHE_SIZE:-0}" -gt 512 ]; then
            log "Shader cache size: ${CACHE_SIZE}MB, cleaning..."
            clear_gpu_cache
        fi
    fi
done
