# root/etc/s6-overlay/s6-rc.d/svc-orca-watchdog/run
#!/usr/bin/with-contenv bash

# OrcaSlicer Watchdog Service
# Monitors OrcaSlicer for hangs, crashes and memory issues, performs automatic recovery

WATCHDOG_INTERVAL=30  # Check every 30 seconds
MEMORY_THRESHOLD=80   # Memory usage percentage threshold
CPU_HANG_THRESHOLD=5  # Minutes of 100% CPU = potential hang
RESTART_COOLDOWN=15   # Seconds to wait after restart before monitoring again
MAX_RESTARTS=10       # Max restarts before giving up (resets after stable period)
STABLE_PERIOD=300     # Seconds of stability before resetting restart counter

log() {
    echo "[$(date '+%Y-%m-%d %H:%M:%S')] WATCHDOG: $1"
}

get_orca_pid() {
    pgrep -f "orca-slicer|OrcaSlicer|orca_slicer" | head -1
}

get_memory_usage() {
    local pid=$1
    if [ -n "$pid" ] && [ -f "/proc/$pid/status" ]; then
        local vmrss=$(grep VmRSS /proc/$pid/status 2>/dev/null | awk '{print $2}')
        echo "${vmrss:-0}"
    else
        echo "0"
    fi
}

check_process_responsive() {
    local pid=$1
    if [ -n "$pid" ]; then
        kill -0 "$pid" 2>/dev/null
        return $?
    fi
    return 1
}

clear_gpu_cache() {
    log "Clearing GPU caches..."
    rm -rf /config/.cache/mesa_shader_cache/* 2>/dev/null
    rm -rf /tmp/mesa* 2>/dev/null
    sync
}

start_orca() {
    log "Starting OrcaSlicer..."

    # Clean up stale lock file and leftover processes
    rm -f /tmp/orcaslicer.lock
    pkill -9 -f WebKitWebProcess 2>/dev/null
    pkill -9 -f WebKitNetworkProcess 2>/dev/null
    sleep 2

    # Launch via the wrapper script using setsid for proper daemonization
    # The wrapper sources gpu-env.sh and sets all required env vars
    setsid /usr/local/bin/orcaslicer-wrapper </dev/null >/dev/null 2>&1 &

    sleep "$RESTART_COOLDOWN"
    local new_pid=$(get_orca_pid)
    if [ -n "$new_pid" ]; then
        log "OrcaSlicer started successfully (PID $new_pid)"
    else
        log "WARNING: OrcaSlicer failed to start"
    fi
}

restart_orca() {
    log "Initiating OrcaSlicer restart..."
    local pid=$(get_orca_pid)

    if [ -n "$pid" ]; then
        log "Sending SIGTERM to PID $pid"
        kill -TERM "$pid" 2>/dev/null
        sleep 5
        if kill -0 "$pid" 2>/dev/null; then
            log "Force killing PID $pid"
            kill -9 "$pid" 2>/dev/null
        fi
    fi

    clear_gpu_cache
    sync
    echo 3 > /proc/sys/vm/drop_caches 2>/dev/null || true

    start_orca
}

log "Starting OrcaSlicer watchdog service"
log "Memory threshold: ${MEMORY_THRESHOLD}%"
log "Check interval: ${WATCHDOG_INTERVAL}s"

HIGH_CPU_COUNT=0
RESTART_COUNT=0
LAST_RESTART_TIME=0

# Wait for display server to be ready
sleep 10

while true; do
    sleep "$WATCHDOG_INTERVAL"

    ORCA_PID=$(get_orca_pid)
    NOW=$(date +%s)

    if [ -z "$ORCA_PID" ]; then
        HIGH_CPU_COUNT=0

        if [ "$RESTART_COUNT" -ge "$MAX_RESTARTS" ]; then
            log "Max restarts ($MAX_RESTARTS) reached. Waiting 5 minutes before trying again."
            sleep 300
            RESTART_COUNT=0
        fi

        RESTART_COUNT=$((RESTART_COUNT + 1))
        log "OrcaSlicer not running (restart $RESTART_COUNT/$MAX_RESTARTS)"
        start_orca
        LAST_RESTART_TIME=$NOW
        continue
    fi

    # Process is alive - track stability
    if [ "$LAST_RESTART_TIME" -gt 0 ] && [ $((NOW - LAST_RESTART_TIME)) -gt "$STABLE_PERIOD" ]; then
        if [ "$RESTART_COUNT" -gt 0 ]; then
            log "OrcaSlicer stable for ${STABLE_PERIOD}s, resetting restart counter"
            RESTART_COUNT=0
        fi
    fi

    # Check if process is responsive
    if ! check_process_responsive "$ORCA_PID"; then
        log "OrcaSlicer (PID $ORCA_PID) is not responsive"
        restart_orca
        HIGH_CPU_COUNT=0
        continue
    fi

    # Check memory usage (in KB)
    MEM_USAGE=$(get_memory_usage "$ORCA_PID")
    MEM_USAGE_MB=$((MEM_USAGE / 1024))
    TOTAL_MEM=$(grep MemTotal /proc/meminfo | awk '{print $2}')
    MEM_PERCENT=$((MEM_USAGE * 100 / TOTAL_MEM))

    if [ "$MEM_PERCENT" -gt "$MEMORY_THRESHOLD" ]; then
        log "High memory usage detected: ${MEM_USAGE_MB}MB (${MEM_PERCENT}%)"
        clear_gpu_cache
        sleep 5
        NEW_MEM=$(get_memory_usage "$ORCA_PID")
        NEW_MEM_PERCENT=$((NEW_MEM * 100 / TOTAL_MEM))
        if [ "$NEW_MEM_PERCENT" -gt "$MEMORY_THRESHOLD" ]; then
            log "Memory still high after cache clear, restarting OrcaSlicer"
            restart_orca
        fi
        HIGH_CPU_COUNT=0
        continue
    fi

    # Check for CPU hang (stuck at 100%)
    CPU_USAGE=$(ps -p "$ORCA_PID" -o %cpu= 2>/dev/null | cut -d. -f1)
    if [ "${CPU_USAGE:-0}" -gt 95 ]; then
        HIGH_CPU_COUNT=$((HIGH_CPU_COUNT + 1))
        log "High CPU detected: ${CPU_USAGE}% (count: $HIGH_CPU_COUNT)"
        if [ "$HIGH_CPU_COUNT" -ge $((CPU_HANG_THRESHOLD * 60 / WATCHDOG_INTERVAL)) ]; then
            log "CPU hang detected, restarting OrcaSlicer"
            restart_orca
            HIGH_CPU_COUNT=0
        fi
    else
        HIGH_CPU_COUNT=0
    fi

    # Periodic cache maintenance (every 10 minutes)
    if [ $((SECONDS % 600)) -lt "$WATCHDOG_INTERVAL" ]; then
        CACHE_SIZE=$(du -sm /config/.cache/mesa_shader_cache 2>/dev/null | cut -f1)
        if [ "${CACHE_SIZE:-0}" -gt 512 ]; then
            log "Shader cache size: ${CACHE_SIZE}MB, cleaning..."
            clear_gpu_cache
        fi
    fi
done
