# root/etc/s6-overlay/s6-rc.d/svc-autocutsel-clipboard/run
#!/usr/bin/with-contenv bash

export DISPLAY="${DISPLAY:-:1}"

# Wait until the X server is reachable so autocutsel can grab the selection.
for i in $(seq 1 60); do
    if xdpyinfo -display "$DISPLAY" >/dev/null 2>&1; then
        break
    fi
    sleep 1
done

exec autocutsel -selection CLIPBOARD -display "$DISPLAY"
